# ipython/notebook.conf
[program:notebook]
directory=/ipy
command=/start.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true